# supervisor.conf
[program:pycursos]
command=/home/ubuntu/pycursos/gunicorn_start
directory=/home/ubuntu/pycursos
user=ubuntu
stdout_logfile=/home/ubuntu/logs/gunicorn_supervisor.log
redirect_stderr=true
